# trunk/target/linux/x64/edge500/files/etc/init.d/vcmakedumpfile
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

START=01

COREDUMP_TOP_DIR=/velocloud
COREDUMP_PATH=/velocloud/kcore
CURR_DATE=`date +%Y-%m-%d:%H:%M:%S`

boot() {
        if [ -e /proc/vmcore ] ; then
		/bin/mkdir -p ${COREDUMP_TOP_DIR}
		/sbin/block mount
                /bin/mkdir -p ${COREDUMP_PATH}
		/sbin/makedumpfile --dump-dmesg /proc/vmcore ${COREDUMP_PATH}/kernel-crash-dmesg.${CURR_DATE}
		/bin/touch ${COREDUMP_PATH}/.kernel-crash-dmesg.${CURR_DATE}
		sync
                /sbin/reboot -f
        fi
}
