# trunk/target/linux/x64/edge500/files/etc/init.d/vcprocesskdump
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

START=99

start() {

	mount /dev/sda6 /velocloud

	[ -e /velocloud/core/kernel_dump ] && {
		/root/kdump_package
	}
}
