# trunk/package/mfgtest/target/bladerf@.service
[Unit]
Description=Run BladeRF
#BindsTo=dev=%i.device

[Service]
Type=simple
ExecStart=/home/velocloud/bladerf/udev.sh add %I
#ExecStop=/home/velocloud/bladerf/udev.sh remove %I
