# trunk/target/linux/x64/edge5x0/files/etc/init.d/vckdumpload
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

START=94

FIND_SUBSTR="root="

start() {

    if [ -e /usr/sbin/kexec ] ; then

		cmdline=`cat /proc/cmdline`

		for arg in $cmdline
		do
			strstr=`expr "$arg" : "$FIND_SUBSTR"`

			if [ ${#FIND_SUBSTR} -eq $strstr ] ; then
				/usr/sbin/kexec -p /vmlinuz.kdump --command-line="1 debug init=/sbin/init earlyprintk=serial,ttyS1,115200 console=ttyS1,115200n8 irqpoll reset_devices earlyprintk=dbgp rootwait maxcpus=1 usbcore.old_scheme_first=0 reboot=p,c $arg"
				break
		        fi
		done
    fi
}
