# trunk/target/linux/x64/edge1000/files/etc/uci-defaults/20-network
#!/bin/sh

# base-files sticks in an /etc/config/network that is useless for us.
# we depend on this file being removed after first boot, so that we
# only clobber the default /etc/config/network on the very first boot.
# [ -r /etc/config/network ] && exit 0

if ! fgrep -q '### DUMMY CONFIGURATION' /etc/config/network 2>/dev/null; then
	# Not the dummy config. All is OK; just exit and remove this script.
	exit 1
fi

EDGETYPE=`/etc/edgetype`

if fgrep -q "### DUMMY CONFIGURATION $EDGETYPE" /etc/config/network 2>/dev/null; then
	# Already the default configuration for this edge type
	exit 1
fi

echo "### DUMMY CONFIGURATION $EDGETYPE" > /etc/config/network
case "$EDGETYPE" in
    edge1000ifm)
        cat /etc/config-defaults/edge1000ifm/network >> /etc/config/network
        ;;
    edge800sm)
        cat /etc/config-defaults/edge800sm/network >> /etc/config/network
        ;;
    edge1000sm-sfp)
        cat /etc/config-defaults/edge1000sm-sfp/network >> /etc/config/network
        ;;
    edge1000sm-qat)
        cat /etc/config-defaults/edge1000sm-qat/network >> /etc/config/network
        ;;
    edge1000sm-gbe)
        cat /etc/config-defaults/edge1000sm-gbe/network >> /etc/config/network
        ;;
    *)  # something else?
        cat /etc/config-defaults/edge1000ifm/network >> /etc/config/network
        ;;
esac
exit 1
