# trunk/feeds/luci/contrib/package/freifunk-policyrouting/files/etc/hotplug.d/.svn/entries
10

dir
10274
http://svn.luci.subsignal.org/luci/trunk/contrib/package/freifunk-policyrouting/files/etc/hotplug.d
http://svn.luci.subsignal.org



2013-06-02T23:53:23.243575Z
9836
soma














ab181a69-ba2e-0410-a84d-ff88ab4c47bc

iface
dir

