# trunk/target/linux/x64/vc-xen-aws/files/etc/init.d/vcearlyinit
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

# early-boot actions for VC edge

START=05

# turn all RGB LEDs off;
#echo 0 > /sys/devices/platform/vc/leds/off

# On Hyper-V, make sure that the Hyper-V drivers are loaded
PRODUCT_NAME=`cat /sys/class/dmi/id/product_name`
SYS_VENDOR=`cat /sys/class/dmi/id/sys_vendor`
if [ "$SYS_VENDOR" = "Microsoft Corporation" -a "$PRODUCT_NAME" = "Virtual Machine" ]; then
    modprobe hv_balloon
    modprobe hv_netvsc
    modprobe hv_utils
fi

VC_PART=`block info | fgrep 'LABEL="user"' | sed -e 's/:.*$//'`
if [ -b "$VC_PART" ]; then
    echo "" > /dev/kmsg
    echo "FIRST TIME: EXPANDING /velocloud, TAKES A FEW MINUTES ..." > /dev/kmsg
    resize2fs "$VC_PART"
    echo "... DONE" > /dev/kmsg
    echo "" > /dev/kmsg
fi
