# trunk/target/linux/x64/edge500/files/etc/uci-defaults/20-network
#!/bin/sh

# base-files sticks in an /etc/config/network that is useless for us.
# we depend on this file being removed after first boot, so that we
# only clobber the default /etc/config/network on the very first boot.
# [ -r /etc/config/network ] && exit 0

if ! fgrep -q '### DUMMY CONFIGURATION' /etc/config/network 2>/dev/null; then
	# Not the dummy config. All is OK; just exit and remove this script.
	exit 0
fi

EDGETYPE=`/etc/edgetype`

case "$EDGETYPE" in
    edge300)
        cp /etc/config-defaults/edge300/network /etc/config/network
        ;;
    edge800p1)
        cp /etc/config-defaults/edge800p1/network /etc/config/network
        ;;
    *)  # edge500 or anything else
        cp /etc/config-defaults/edge500/network /etc/config/network
        cp /etc/config-defaults/edge500/wireless /etc/config/wireless
        ;;
esac
