# trunk/feeds/packages/net/madwifi/files/etc/hotplug.d/.svn/entries
10

dir
40995
svn://svn.mein.io/openwrt/packages/net/madwifi/files/etc/hotplug.d
svn://svn.mein.io/openwrt



2008-07-30T05:22:16.679739Z
12005
nbd














3c298f89-4303-0410-b956-a3cf2f4a3e73

net
dir

