# trunk/feeds/luci/contrib/package/freifunk-firewall/files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 79
/!svn/ver/9832/luci/trunk/contrib/package/freifunk-firewall/files/etc/hotplug.d
END
